# scripts/inventory-system.service
[Unit]
Description=Update mini-pc-inventory repo
After=update-inventory-script.target

[Service]
WorkingDirectory=/srv/mini-pc-inventory/scripts
ExecStart=/srv/mini-pc-inventory/scripts/runit.sh

[Install]
WantedBy=multi-user.target
